# examples/cron.d/backup-vms
# m h  dom mon dow   command
00 04 * * 6 root /usr/local/sbin/backup-vms -c -f /etc/backup-vms/conf.d/k7-dispatch.conf
00 04 * * 7 root /usr/local/sbin/backup-vms -c -f /etc/backup-vms/conf.d/all-others.conf
